# .config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver Service

[Service]
ExecStart=/usr/bin/xscreensaver

[Install]
WantedBy=default.target
